# Operations/panic_switch/turn_back_win.service
[Unit]
Description=Time-Limited Turn_back_Win Service
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/turn_back_win.sh
StandardOutput=journal
StandardError=journal
# Ensure TimeoutStopSec is generous enough for the total duration + buffer
TimeoutStopSec=15min

[Install]
WantedBy=multi-user.target